# .docker/prod/crontab
# 定期実行の設定
# ここを編集して優待速報監視のタイミングを変える

# 例
# 平日 8-17時 → 10分おき
*/10 8-17 * * 1-5 /app/yutai_watch

# それ以外 → 30分おき
*/30 0-7,18-23 * * * /app/yutai_watch
*/30 * * * 0,6 /app/yutai_watch