# .docker/prod/crontab
# 定期実行の設定
# ここを編集して優待速報監視のタイミングを変える

# Rustのジッターに合わせるため、なんとなく位相を1分早めるめんどくさい書き方をしている。
# 本来やりたい設定
# 平日(月-金) 8時～15時半 → 10分
# 平日(月-金) 15時半～24時 → 30分
# 平日(月-金) 24時～8時 なし
# 休日(土-日) なし

# ジッターを考慮した設定
# 平日(月-金) 7時59分～15時29分 → 10分
# 平日(月-金) 15時59分～23時59分 → 30分
# 平日(月-金) 24時～8時 なし
# 休日(土-日) なし

# =========================
# monitor 実行
# =========================
# -------------------------
# 平日：メイン稼働（7:59〜15:29）10分おき
# -------------------------
# 7時台だけ例外（スタート調整）
59 7 * * 1-5 /app/yutai_watch -- monitor

# 9時〜14時：10分おき
9,19,29,39,49,59 8-14 * * 1-5 /app/yutai_watch -- monitor

# 15時：15:29まで
9,19,29 15 * * 1-5 /app/yutai_watch -- monitor

# -------------------------
# 平日：それ以外の時間（15:29〜23:59）30分おき
# -------------------------
# 夜間（平日）
59 15 * * 1-5 /app/yutai_watch -- monitor
29,59 16-23 * * 1-5 /app/yutai_watch -- monitor

# -------------------------
# 平日：それ以外の時間（翌24:00〜翌8:00）なし
# -------------------------

# -------------------------
# 休日（土日）：終日なし
# -------------------------


# =========================
# prune 実行
# =========================
# 毎月1日の10:00
0 10 1 * * /app/yutai_watch -- prune



# ============================================================
# cron 設定メモ
# ============================================================
# cron の基本構文
#
# ┌ 分 (0-59)
# │ ┌ 時 (0-23)
# │ │ ┌ 日 (1-31)
# │ │ │ ┌ 月 (1-12)
# │ │ │ │ ┌ 曜日 (0-6) ※ 0=日曜
# │ │ │ │ │
# * * * * * 実行コマンド
#
# ------------------------------------------------------------
# 曜日
# ------------------------------------------------------------
# 0 = 日曜日
# 1 = 月曜日
# 2 = 火曜日
# 3 = 水曜日
# 4 = 木曜日
# 5 = 金曜日
# 6 = 土曜日
#
# 例:
# 平日だけ → 1-5
# 土日だけ → 0,6
#
# ------------------------------------------------------------
# よく使う書き方
# ------------------------------------------------------------
# *     = 毎回
# */10  = 10分おき
# */30  = 30分おき
# */2   = 2時間おき
#
# 1,5,10 = 指定値のみ
# 1-5    = 範囲指定
#
# ------------------------------------------------------------
# 実行例
# ------------------------------------------------------------
# 毎日0時
# 0 0 * * * command
#
# 毎時30分
# 30 * * * * command
#
# 平日9-17時の10分おき
# */10 9-17 * * 1-5 command
#
# 毎月1日の10時
# 0 10 1 * * command
#
# ------------------------------------------------------------
# コメントアウト
# ------------------------------------------------------------
# 行頭に # を付けると無効化される
#
# 例:
# # 0 * * * * command
#
# ------------------------------------------------------------
# 注意点
# ------------------------------------------------------------
#
# cron は「15:30まで」のような細かい終了条件を
# 1行で書けない場合がある
#
# そのため必要に応じて
# 「15時台だけ別指定」
# のように分けて書く